# roles/admin-notifier/files/etc/cron.d/alerta-heartbeat
PATH=/bin:/usr/bin:/sbin:/usr/sbin
* * * * * root /usr/local/bin/alerta-notify heartbeat $(hostname) "" ""
