# roles/admin-notifier/files/etc/systemd/system/admin-notifier-battery.service
[Unit]
Description=Check battery and inform alerta if low

[Service]
Type=oneshot
ExecStart=/etc/admin-notifier/battery

[Install]
WantedBy=multi-user.target
