# roles/admin-notifier/files/etc/systemd/system/admin-notifier-traefik.timer
[Unit]
Description=Check traefik status and inform alerta if down every 10min

[Timer]
OnCalendar=*-*-* *:0/10:00

[Install]
WantedBy=timers.target
