# roles/admin-notifier/files/etc/systemd/system/admin-notifier-ping.service
[Unit]
Description=Send ping to alerta

[Service]
Type=oneshot
ExecStart=/etc/admin-notifier/ping

[Install]
WantedBy=multi-user.target
