# roles/admin-notifier/files/etc/systemd/system/admin-notifier-battery.timer
[Unit]
Description=Check battery and inform alerta if low every 15min

[Timer]
OnCalendar=*-*-* *:0/15:00

[Install]
WantedBy=timers.target
