# roles/admin-notifier/files/etc/systemd/system/admin-notifier-ping.timer
[Unit]
Description=Send ping to alerta every minutes

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
