# roles/admin-notifier/files/etc/systemd/system/admin-notifier-apt.timer
[Unit]
Description=Check new packages and inform alerta every day at 10am.

[Timer]
OnCalendar=*-*-* 10:00:00

[Install]
WantedBy=timers.target
