# roles/admin-notifier/files/etc/systemd/system/admin-notifier-systemd.service
[Unit]
Description=Check failed services and inform alerta if any

[Service]
Type=oneshot
ExecStart=/etc/admin-notifier/systemd

[Install]
WantedBy=multi-user.target
